# nixpkgs/xdg/s6u/sv/shepherd-srv/run
#!/bin/execlineb -P
envfile /etc/s6/config/user-services.conf
importas XDG_RUNTIME_DIR XDG_RUNTIME_DIR
exec s6-fghack /usr/local/bin/shepherd
