# s6u/sv/birdtray/run
#!/bin/execlineb -P
s6-envdir ../../envdir

fdmove -c 2 1 birdtray
