# s6-rc/sv/pw-pulse-srv/run
#!/bin/execlineb -P
#redirfd -w 1 /dev/null pipewire -c pipewire-pulse.conf
backtick -E USER { id -un } s6-envdir /home/${USER}/.config/s6u/envdir

fdmove -c 2 1 pipewire -c pipewire-pulse.conf
