# nixpkgs/xdg/s6u/sv/wireplumber/run
#!/bin/execlineb -P
#redirfd -w 1 /dev/null wireplumber
envfile /etc/s6/config/user-services.conf

fdmove -c 2 1 wireplumber
