# s6u/sv/udiskie-srv/run
#!/bin/execlineb -P
s6-envdir ../../envdir
importas HOME HOME
fdmove -c 2 1  udiskie --tray
