# s6u/sv/swhkd-srv/run
#!/bin/execlineb -P
envfile /etc/s6/config/user-services.conf
fdmove -c 2 1 pkexec /usr/bin/swhkd
