# s6-rc/sv/rsync-calibre/up
#!/bin/execlineb -P
backtick -E USER { id -un } s6-envdir /home/${USER}/.config/s6-rc/envdir
importas HOME HOME
fdmove -c 2 1 rsync -avz --progress ${HOME}/Calibre-Library/ vector-sigma:/ssd-128/Calibre-library
