# s6u/sv/sxhkd-srv/run
#!/bin/execlineb -P
s6-envdir ../../envdir
importas HOME HOME
importas DISPLAY DISPLAY
fdmove -c 2 1  sxhkd
