# s6u/sv/dunst-srv/run
#!/bin/execlineb -P
s6-envdir ../../envdir
importas DISPLAY DISPLAY
fdmove -c 2 1 dunst
