# s6-rc/sv/pw-log/run
#!/bin/execlineb -P
importas USER USER
fdmove -c 2 1 /home/${USER}/.config/s6-rc/user-logger pipewire
