# s6u/sv/pw-srv/run
#!/bin/execlineb -P
#redirfd -w 1 /dev/null pipewire
s6-envdir ../../envdir

fdmove -c 2 1 pipewire
