# s6u/sv/pw-pulse-srv/run
#!/bin/execlineb -P
#redirfd -w 1 /dev/null pipewire -c pipewire-pulse.conf
envfile /etc/s6/config/user-services.conf

fdmove -c 2 1 pipewire -c pipewire-pulse.conf
