# s6-rc/sv/dbus-srv/run
#!/bin/execlineb -P
backtick -E USER { id -un } s6-envdir /home/${USER}/.config/s6-rc/envdir
importas -i ADDR DBUS_SESSION_BUS_ADDRESS
fdmove -c 2 1
dbus-daemon --nofork --nopidfile --session --address=${ADDR} --print-address=3
