# s6u/sv/emacs-exwm-srv/run
#!/bin/execlineb -P
s6-envdir ../../envdir
importas HOME HOME
fdmove -c 2 1  emacs --chdir=${HOME} --fg-daemon -f exwm-enable
