# s6-rc/sv/ssh-agent-srv/run
#!/bin/execlineb -P
backtick -E USER { id -un } s6-envdir /home/${USER}/.config/s6-rc/envdir
importas XDG_RUNTIME_DIR XDG_RUNTIME_DIR
fdmove -c 2 1  /usr/bin/ssh-agent -D -a ${XDG_RUNTIME_DIR}/ssh-agent.socket
