# s6u/sv/herbe/run
#!/bin/execlineb -P
         s6-envdir ../../envdir
         importas DISPLAY DISPLAY
         pipeline { jq --unbuffered ".summary,.body" }
                  forstdin-E  $notfi
                    herbe $noti
