# s6-rc/sv/keynav/run
#!/bin/execlineb -P
backtick -E USER { id -un }
s6-envdir /home/${USER}/.config/s6u/envdir/
importas DISPLAY DISPLAY
fdmove -c 2 1 keynav
