# s6u/sv/tiramisu/run
#!/bin/execlineb -P
redirfd -wnb 1 /tmp/noti-fifo tiramisu
