# nixpkgs/xdg/s6u/sv/lf/run
#!/bin/execlineb -P
fd-move -c 2 1 lf -server
