# nixpkgs/xdg/s6u/sv/gammastep-log/run
#!/bin/execlineb -P
backtick -E USER { id -un }
importas USER USER
fdmove -c 2 1 /home/${USER}/.config/s6u/user-logger gammastep
